# build/nzbget.conf
[supervisord]
loglevel = info

[program:container-log-forwarder]
autorestart = true
startsecs = 0
user = root
priority = 1
command = /bin/bash -c "touch /data/nzbgetvpn-container.log && tail -n 0 -F /data/nzbgetvpn-container.log"
umask = 000
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0

[program:start-script]
autorestart = false
startsecs = 0
user = root
priority = 10
command = /root/start.sh
umask = 000
stdout_logfile = /data/nzbgetvpn-container.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 5
redirect_stderr = true

[program:watchdog-script]
autorestart = false
startsecs = 0
user = nobody
priority = 10
command = /home/nobody/watchdog.sh
umask = 000
stdout_logfile = /data/nzbgetvpn-container.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 5
redirect_stderr = true

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
priority = 10
command = /usr/local/bin/shutdown.sh 'nzbget'
umask = 000
stdout_logfile = /data/nzbgetvpn-container.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 5
redirect_stderr = true
