# build/nzbget.conf
[supervisord]
loglevel = info

[program:start-script]
autorestart = false
startsecs = 0
user = root
command = /root/start.sh
umask = 000
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0

[program:watchdog-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/watchdog.sh
umask = 000
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh 'nzbget'
umask = 000
stdout_logfile = /dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/fd/2
stderr_logfile_maxbytes = 0
